# deploy/shortnames.conf
description "Short Names"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec /usr/src/short_names/deploy/gunicorn.sh
